# ansible/roles/raspicam-timelapse/templates/raspicam-timelapse.service
# /etc/systemd/system/raspicam-timelapse.service
[Unit]
Description=raspicam-timelapse service

[Service]
User=pi
SyslogIdentifier=raspicam-timelapse
SyslogFacility=user
LimitNOFILE=65535
LimitCORE=infinity
Environment=NODE_ENV=production
WorkingDirectory=/home/pi/raspicam-timelapse
ExecStart={{node_path}}/node server.js
Restart=always

[Install]
WantedBy=multi-user.target
